# config/systemd/ibkrdata-backup.service
[Unit]
Description=IBKR Data Database Backup

[Service]
Type=oneshot
WorkingDirectory=%h/projects/ibkrData
ExecStart=%h/projects/ibkrData/db/backup.sh --cron
